# server/ansible/roles/preview-manager/templates/preview-manager.service.j2
[Unit]
Description=Preview Manager Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=preview-user
Group=preview-user
WorkingDirectory=/home/preview-manager/www/previews/server/preview-manager
Environment="PATH=/home/preview-manager/www/previews/server/preview-manager/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/preview-manager/www/previews/server/preview-manager/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=preview-manager

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
# MemoryLimit=2G
# CPUQuota=200%

[Install]
WantedBy=multi-user.target
